# ansible/image/templates/mw_img_backup.timer.j2
[Unit]

Description=MediaWiki image folder backup service timer

[Timer]

OnCalendar={{ backup_timer_setting }}

[Install]

WantedBy=timers.target
