# ansible/db/templates/db_dump_backup.timer.j2
[Unit]

Description=MediaWiki database dump backup service timer

[Timer]

OnCalendar={{ backup_timer_setting }}

[Install]

WantedBy=timers.target
